# supervisor/varnishd.conf
[program:varnishd]
command=/bin/bash -c 'source /etc/default/varnish && /usr/sbin/varnishd -F $DAEMON_OPTS'
autorestart=true